# resources/supervisord.conf
[supervisord]
nodaemon=true

[program:xld]
command=/opt/xld/server/bin/run.sh
